# crontab
*/15 * * * * python3 /home/pi/apps/mitemp/mitemp-mqtt.py 58:2D:34:36:82:ED --measurements --health